# sci-libs/ni-pal/files/ni-pal-2.5.4_beta2-nipalwrapper
#!/sbin/runscript

depend()
{
   provide nipalsm
}

start()
{
   /etc/init.d/nipal start
}

stop()
{
   /etc/init.d/nipal stop
}

