# sci-libs/ni-gpib/files/ni-gpib-2.5_beta-gpibenumsvcwrapper
#!/sbin/runscript

depend()
{
   need nipalsm
}

start()
{
   /etc/init.d/gpibenumsvc start
}

stop()
{
   /etc/init.d/gpibenumsvc stop
}
